# agent-ui.service
[Unit]
Description=Agent UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/my/agent-ui
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/make run
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-ui

[Install]
WantedBy=multi-user.target
